# roles/prepare/templates/retrospring-sidekiq.service.j2
# {{ ansible_managed }}
# container-retrospring-sidekiq.service
# autogenerated by Podman 4.3.1
# Thu Feb  2 19:29:17 UTC 2023

[Unit]
Description=Podman container-retrospring-sidekiq.service
Documentation=man:podman-generate-systemd(1)
Wants=network-online.target
After=network-online.target
RequiresMountsFor=%t/containers

[Service]
Environment=PODMAN_SYSTEMD_UNIT=%n
Restart=on-failure
TimeoutStopSec=70
ExecStartPre=/bin/rm \
    -f %t/%n.ctr-id
ExecStart=/usr/bin/podman container run \
    --cidfile=%t/%n.ctr-id \
    --cgroups=no-conmon \
    --rm \
    --sdnotify=conmon \
    -d \
    --replace \
    --name retrospring-sidekiq \
    -v /home/justask/devise.rb:/opt/retrospring/app/config/initializers/devise.rb \
    -v /home/justask/justask.yml:/opt/retrospring/app/config/justask.yml \
    -v /home/justask/puma.rb:/opt/retrospring/app/config/puma.rb \
    -v /home/justask/sidekiq.yml:/opt/retrospring/app/config/sidekiq.yml \
{% for item in _env_vars | dict2items %}
    -e {{ (item.key + "=" + item.value) | quote }} \
{% endfor %}
    -p 3000:3000 docker.io/retrospring/retrospring:{{ rs_version }} bundle exec sidekiq
ExecStop=/usr/bin/podman stop \
    --ignore -t 10 \
    --cidfile=%t/%n.ctr-id
ExecStopPost=/usr/bin/podman rm \
    -f \
    --ignore -t 10 \
    --cidfile=%t/%n.ctr-id
Type=notify
NotifyAccess=all

[Install]
WantedBy=default.target
